# deploy/goamines.service
[Unit]
Description=Datasette - Goa Ore Transport (goamines)
After=network.target

[Service]
Type=simple
User=goamines
Group=goamines
WorkingDirectory=/home/goamines/goamines
# Serve immutable (-i): the DB is built offline and never written in place, so
# SQLite needs no WAL/lock files and Datasette can cache aggressively.
# Datasette 1.0: descriptive metadata in metadata.yaml, config (settings, plugins,
# canned queries) in datasette.yaml (-c). Settings live in datasette.yaml, not CLI flags.
ExecStart=/home/goamines/goamines/.venv/bin/datasette serve -i goamines.db \
  -m metadata.yaml -c datasette.yaml --static static:static \
  -h 127.0.0.1 -p 8001
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
